# var/lib/dpkg/info/iot.list
/.
/etc
/etc/init.d
/etc/init.d/iot
/etc/iotsample-raspberrypi
/opt
/opt/iot
/opt/iot/IoTFoundation.pem
/opt/iot/iot
/opt/iot/iotGetDeviceID.sh
/opt/iot/README
/usr
/usr/lib
/usr/lib/libpaho-mqtt3a.so.1.0
/usr/lib/libpaho-mqtt3as.so.1.0
/usr/share
/usr/share/doc
/usr/share/doc/iot
/usr/share/doc/iot/changelog.Debian.gz
/usr/share/doc/iot/README.Debian.gz
/usr/share/doc/iot/copyright
/usr/lib/libpaho-mqtt3as.so
/usr/lib/libpaho-mqtt3a.so.1
/usr/lib/libpaho-mqtt3as.so.1
/usr/lib/libpaho-mqtt3a.so
